# eumetsat-wallpaper.service
[Unit]
Description=EUMETSAT Linux wallpaper updater (example)
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/eumetsat-wallpaper run
WorkingDirectory=%h
Environment=PATH=%h/.local/bin:%h/bin:/usr/local/bin:/usr/bin

[Install]
WantedBy=default.target
